# launchd/com.aihugh.morning.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.aihugh.morning</string>

  <!-- Push the morning briefing notification at 06:30 local time, weekdays. -->
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>30</integer></dict>
  </array>

  <key>WorkingDirectory</key>
  <string>~/Developer/ai-hugh</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>~/Developer/ai-hugh/morning_notify.sh</string>
  </array>

  <key>StandardOutPath</key>
  <string>~/Developer/ai-hugh/logs/morning.out.log</string>
  <key>StandardErrorPath</key>
  <string>~/Developer/ai-hugh/logs/morning.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
